# systemd/gametimescheduler.service
[Unit]
Description=Discord Scheduler Bot
After=network.target

[Service]
WorkingDirectory=/home/pi/python-discord-scheduler-bot
EnvironmentFile=/etc/environment
ExecStart=/home/pi/python-discord-scheduler-bot/.venv/bin/python3 /home/pi/python-discord-scheduler-bot/bot.py
Restart=on-failure
User=pi
Group=pi

[Install]
WantedBy=multi-user.target